# systemd/admin-service.service
[Unit]
Description=Admin Service - Auth System Management Interface
Documentation=https://github.com/your-org/um_oic
After=network.target auth-service.service
Requires=auth-service.service

[Service]
Type=simple
User=auth
Group=auth
WorkingDirectory=/var/lib/auth

# Service configuration
ExecStart=/usr/local/bin/admin-service \
  --data-dir /var/lib/auth/data \
  --cert-dir /etc/ssl/auth \
  --config /etc/auth/config.toml \
  --bind 0.0.0.0:8001 \
  --auth-pid-file /var/run/auth-service.pid

# Stop service
ExecStop=/bin/kill -TERM $MAINPID

# Restart policy
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=admin-service

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/auth/data /var/run
ReadOnlyPaths=/var/lib/auth /etc/auth /etc/ssl/auth

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=RUST_LOG=info
Environment=ADMIN_DATA_DIR=/var/lib/auth/data
Environment=ADMIN_CERT_DIR=/etc/ssl/auth
Environment=AUTH_PID_FILE=/var/run/auth-service.pid

[Install]
WantedBy=multi-user.target